# systemd/memory-platform-maintenance.service.template
[Unit]
Description=Memory Platform bounded maintenance
After=network-online.target

[Service]
Type=oneshot
ExecStart=__ROOT__/scripts/run-memory-maintenance.sh --daily
