# systemd/user/memory-archive-verify.timer
[Unit]
Description=Daily memory archive verification

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true

[Install]
WantedBy=timers.target
