# systemd/user/memory-ingest.timer
[Unit]
Description=Run memory freshness ingest daily

[Timer]
OnCalendar=*-*-* 03:25:00
Persistent=true
RandomizedDelaySec=10m
Unit=memory-ingest.service

[Install]
WantedBy=timers.target
