# systemd/user/memory-ingest.service
[Unit]
Description=Memory platform freshness ingest
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/humanoracle26/memory-platform-rust
EnvironmentFile=-/home/humanoracle26/memory-platform-rust/.env
Environment=RUST_LOG=info
ExecStart=/home/humanoracle26/memory-platform-rust/ingest.sh all
